# configs/osagent.conf
description "OS Agent"
author "Vitor Leonel <falecom@vitorleonel.com>"

start on filesystem
stop on shutdown

respawn

chdir /home/osagent
exec node src
